# templates/avahi-aliases.conf.j2
[program:avahi-aliases]
command=/usr/local/bin/avahi-cname-aliases
autostart=true
autorestart=true
user=nobody
directory=/var/run
